# etc/cron.d/minereu
10 *     * * *     root   /usr/sbin/ntpdate -u 0.uk.pool.ntp.org
* *     * * *     root  /opt/scripta/startup/miner-start.sh